# deploy/ai.ryvion.node.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.ryvion.node</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ryvion-node</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RYV_HUB_URL</key>
        <string>https://api.ryvion.ai</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/usr/local/var/ryvion</string>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/ryvion-node.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/ryvion-node.err</string>
</dict>
</plist>
